# wiki-monitor.service
[Unit]
Description=Wiki-OS file monitor daemon
Documentation=file://%h/wiki-os/README.md
# Start after the user's graphical session or basic login session is up.
After=default.target
# If Ollama is managed as a user service, uncomment the next line:
# After=ollama.service

[Service]
Type=simple
# The daemon runs as the login user — no root, no special capabilities.
ExecStart=%h/wiki-os/.venv/bin/python3 -m src.monitor
WorkingDirectory=%h/wiki-os
Environment=WIKI_OS_CONFIG=%h/.config/wiki-os/config.json
Environment=PYTHONPATH=%h/wiki-os

# Restart on failure with a 10-second delay.
# This handles the case where Ollama is not yet running when the monitor starts.
Restart=on-failure
RestartSec=10

# Write stdout/stderr to the systemd journal. Access with:
#   journalctl --user -u wiki-monitor -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
