# systemd/wiki-ssd-backup.service
[Unit]
Description=Wiki SSD Backup On Attach
After=default.target

[Service]
Type=oneshot
ExecStart=%h/Documents/wiki-linux/wiki-linux/bin/wiki-ssd-backup
WorkingDirectory=%h/Documents/wiki-linux/wiki-linux
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
