# systemd/wiki-kb-sync.service
[Unit]
Description=Wiki Knowledge Base Sync to OpenWebUI
After=wiki-openwebui.service

[Service]
Type=oneshot
ExecStart=%h/Documents/wiki-linux/wiki-linux/bin/wiki-kb-sync
StandardOutput=append:%h/.local/share/wiki-linux/logs/wiki-kb-sync.log
StandardError=append:%h/.local/share/wiki-linux/logs/wiki-kb-sync.log
