# systemd/wiki-kb-sync.timer
[Unit]
Description=Wiki Knowledge Base Sync Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=1h
Unit=wiki-kb-sync.service

[Install]
WantedBy=timers.target
